# config/flask_project.conf
[program:flask_project]
command = /home/deploy/env/bin/gunicorn app:app -b localhost:8000
directory = /home/deploy/your_app_name
user = deploy
