# config/flask_project.conf
[program:your_app_name]
command = /home/deploy/.your_app_name/bin/gunicorn app:app -b localhost:8000
directory = /home/deploy/your_app_name
user = deploy
